# files/influxdb.in
#!/bin/sh
#

# PROVIDE: influxdb
# REQUIRE: LOGIN
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf or /etc/rc.conf.d/influxdb
# to configure influxdb.
# influxdb_enable (bool):   Set to "NO" by default.
#           Set it to "YES" to enable influxdb.
# influxdb_config (str):    Default to "%%ETCDIR%%/config.toml"
#           Configuration file.
#

. /etc/rc.subr

name="influxdb"
rcvar=influxdb_enable

load_rc_config $name

: ${influxdb_enable="NO"}
: ${influxdb_config="%%ETCDIR%%/config.toml"}

command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/influxdb"
command_args="-c -f -- ${procname} -config ${influxdb_config}"

run_rc_command "$1"
